# supervisor/lae-cluster.conf
[program:lae-worker]
process_name=%(program_name)s_%(process_num)02d
directory=/opt/lae
command=php artisan cluster:work
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/lae-cluster.log
stopwaitsecs=3600
