# supervisor/lae-default.conf
[program:lae-default]
process_name=%(program_name)s_%(process_num)02d
command=docker run -i --rm --init --restart=always --net=host -v /opt/lae:/opt/lae ccr.ccs.tencentyun.com/laecloud/cafe:latte art queue:work --queue=default
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=root
numprocs=8
redirect_stderr=true
stdout_logfile=/var/log/lae-default.log
stopwaitsecs=3600
